# systemd/remoteshell-agent.service
[Unit]
Description=RemoteShell Agent
Documentation=https://github.com/votre-projet/remoteshell
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/remoteshell
ExecStart=/usr/local/bin/remoteshell-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sécurité
NoNewPrivileges=false
PrivateTmp=false

# Limites de ressources
LimitNOFILE=65536
LimitNPROC=4096

# Variables d'environnement
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Configuration de l'agent (modifiez ces valeurs selon votre environnement)
Environment="REMOTESHELL_SERVER_HOST=10.0.0.59"
Environment="REMOTESHELL_SERVER_PORT=8081"
Environment="REMOTESHELL_AGENT_ID=serveur-impression-01"
Environment="REMOTESHELL_AGENT_NAME=Serveur d'impression principal"
Environment="REMOTESHELL_AUTH_TOKEN=test-token"
# Optionnel : charger depuis un fichier
# EnvironmentFile=/etc/remoteshell/agent.conf

[Install]
WantedBy=multi-user.target

